# configs/systemd/vula-discover.target
[Unit]
Description=vula discover service
