# configs/systemd/vula-publish.target
[Unit]
Description=vula publish service
