# configs/systemd/vula-publish-alt.target
[Unit]
Description=vula publish alt service
