# configs/systemd/vula-organize-monolithic.target
[Unit]
Description=vula monolithic service
