# configs/systemd/vula-discover-alt.target
[Unit]
Description=vula discover alt service
